# etc/rc.local
#!/bin/sh -e
find /etc/boot.d -type f -exec bash {} \;
